# fiscal_wine/files/default/tremol/crontab
*/2 * * * * /opt/knowhowERP/util/fix_ole_error.sh
